# services/embed/systemd/captain-memo-embed.service
[Unit]
Description=Captain Memo embedder sidecar (voyage-4-nano via FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=captain-memo-embed
Group=captain-memo-embed
WorkingDirectory=/opt/captain-memo-embed
Environment="PATH=/opt/captain-memo-embed/venv/bin:/usr/bin:/bin"
Environment="HF_HOME=/opt/captain-memo-embed/models"
ExecStart=/opt/captain-memo-embed/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8124 --log-level info
Restart=on-failure
RestartSec=5
# Hardening — match Aelita's pattern.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/captain-memo-embed/models /opt/captain-memo-embed/logs
ProtectHome=true

[Install]
WantedBy=multi-user.target
